# systemd/college-calendar.service
[Unit]
Description=College Calendar Automation
Documentation=https://github.com/yourusername/college-calendar
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=%u
WorkingDirectory=%h/Projects/college-calender
ExecStart=%h/Projects/college-calender/run_pipeline.sh

# Environment
Environment="PATH=%h/.cargo/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin"

# Logging
StandardOutput=append:%h/Projects/college-calender/logs/service.log
StandardError=append:%h/Projects/college-calender/logs/service-error.log

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=%h/Projects/college-calender

# Timeout
TimeoutStartSec=600

# Restart policy (don't restart on failure for oneshot)
Restart=no

[Install]
WantedBy=default.target
